# setup/git-repo-monitor.service
[Unit]
Description=Abbey Stock Exchange Startup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=orangepi
Group=orangepi
WorkingDirectory=/home/orangepi/AbbeyStockExchange
ExecStart=/home/orangepi/AbbeyStockExchange/start_app.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
# Timeout configuration - service will be killed after 30 seconds
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target