# setup/git-repo-monitor.service
[Unit]
Description=Git Repository Boot Check Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=orangepi
Group=orangepi
ExecStart=/usr/local/bin/git-repo-monitor.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target