# setup/git-repo-monitor.service
[Unit]
Description=Git Repository Boot Check Service
After=graphical-session.target
Wants=network-online.target

[Service]
Type=oneshot
Environment="DISPLAY=:0"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/usr/local/bin/git-repo-monitor.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=default.target